# data/snapper-daily.timer
[Unit]
Description=Daily Snapper snapshot creation/cleanup

[Timer]
OnBootSec=30min
OnUnitActiveSec=1day